# buildroot/board/common-headless/rootfs/service-bootloader.service
[Unit]
Description=Service Buildroot
Wants=multi-user.target
After=network.target multi-user.target

[Service]
ExecStart=/usr/bin/service-bootloader
Restart=always
RestartSec=2